# ansible/roles/endlessh/templates/endlessh.service
[Unit]
Description=Endlessh

[Service]
ExecStart=/usr/local/bin/endlessh -4 -p {{ endlessh.port | default(22) }}